# deploy/systemd/ghost-mail-bridge-api.service
[Unit]
Description=Ghost Mail Bridge API
After=network.target

[Service]
Type=simple
User=ghost-mail-bridge
Group=ghost-mail-bridge
WorkingDirectory=/opt/ghost-mail-bridge/current
EnvironmentFile=/opt/ghost-mail-bridge/config.env
Environment=APP_ROLE=api
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
